# supervisord-phalcon.conf
[program:apache2]
command=/start-phalcon.sh
numprocs=1
autostart=true
autorestart=true
